# hbnb/arm/custom-settings/config/files/root/etc/systemd/system/hbnbuild.service
[Unit]
Description=HotspotBnB Build Final Phase
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/home/pi/hbnbuild.sh -k start -DFOREGROUND

[Install]
WantedBy=default.target
